# client/hostit-agent.service
[Unit]
Description=HostIt Tunnel Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%I
EnvironmentFile=-/etc/hostit/agent.env

# Build only if missing
ExecStartPre=/bin/sh -c 'test -x ./bin/tunnel-agent || ./build.sh'
ExecStart=%I/client.sh

Restart=always
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
