# server/hostit-server.service
[Unit]
Description=HostIt Tunnel Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%I
# Optional environment overrides
EnvironmentFile=-/etc/hostit/server.env

# Build only if missing
ExecStartPre=/bin/sh -c "test -x ./bin/tunnel-server || (echo Missing ./bin/tunnel-server. Run ./build.sh once as your user. >&2; exit 1)"
ExecStart=%I/server.sh


Restart=always
RestartSec=2

# Hardening (light)
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
